# crontab
* * * * * python /root/r6-calendar/main.py >> /root/r6-calendar/cron.log
